# .config/autostart/wayland-xhost-local.desktop
[Desktop Entry]
Name=xhost + local:
Comment=xhost + local: for toolbox X support in SSH
Exec=xhost + local:
Type=Application
